# src/php-env-base/assets/etc/nginx/supervisor/20.program.nginx.ini
[program:nginx]
autostart=true
startsecs=10
stopwaitsecs=5
command=/usr/sbin/nginx -c /etc/nginx/main.nginx -g 'daemon off;'
