# src/phpbrew/assets/template/703-nts/etc/supervisor/10.program.php-fpm.prod.ini

[program:php-fpm.%PHP_VERSION_NAME%.prod]
autostart=false
startsecs=10
command=%(ENV_HOME)s/.phpbrew/php/%PHP_VERSION_NAME%/sbin/php-fpm
	--allow-to-run-as-root
	--nodaemonize
	--fpm-config %(ENV_HOME)s/.phpbrew/php/%PHP_VERSION_NAME%/etc/php-fpm.d/php-fpm.prod.ini
	-c %(ENV_HOME)s/.phpbrew/php/%PHP_VERSION_NAME%/etc/php-fpm.d/php.prod.ini
