# src/php-env-dev/assets/etc/supervisor/conf.d/10.program.chromium.ini
[program:chromium.9222]
autostart = true
command = chromium-browser
	--headless
	--disable-gpu
	--remote-debugging-port=9222
	--user-data-dir=/tmp
	--no-sandbox
	http://localhost
stdout_logfile = /var/log/chromium.headless.9222.stdout.log
stderr_logfile = /var/log/chromium.headless.9222.stderr.log

[program:chromium.driver.4444]
autostart = false
command = /usr/bin/chromedriver
	--port=4444
	--log-level=ALL
	--log-path=/var/log/chromium.driver.4444.log
stdout_logfile = /var/log/chromium.driver.4444.stdout.log
stderr_logfile = /var/log/chromium.driver.4444.stderr.log
